# deploy/tesseract.service
# Systemd service for Tesseract
# Install:
#   1. Update User, WorkingDirectory, ExecStart, and EnvironmentFile
#      to match your server's checkout path and user.
#   2. sudo cp deploy/tesseract.service /etc/systemd/system/
#      sudo systemctl daemon-reload
#      sudo systemctl enable --now tesseract

[Unit]
Description=Tesseract (Rails)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/tesseract
ExecStart=/opt/tesseract/bin/prod
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/tesseract/.env

# Logging
StandardOutput=append:/opt/tesseract/log/production.log
StandardError=append:/opt/tesseract/log/production.log

[Install]
WantedBy=multi-user.target
